# backend/lensing-backend.conf
start on runlevel 2
env LOC=/var/www/html/checkouts/20150115/domains/wikipediahistory/backend
exec python2.7 $LOC/backend -p 1530 -c $LOC/full.conf >> upstart.log 2>&1
respawn
