# backend/lensingup.conf
description "lensingwikipedia"
start on runlevel 2
exec $BACKEND/backend -c $BACKEND/full.conf
respawn
respawn limit 10 3600
